# packaging/systemd/btrfs-manager-snapshot@.timer
[Unit]
Description=Btrfs Manager scheduled snapshot policy for %i

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
